# supervisor-tobacco-api.ini
[program:tobacco-api]
directory=/mnt/tobacco/TobaccoOA
command = /mnt/tobacco/TobaccoOA/venv/bin/uwsgi --ini /mnt/tobacco/TobaccoOA/uwsgi.ini
stderr_logfile = /mnt/tobacco/TobaccoOA/run.log
